# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/local/apache2/bin/httpd -k start
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D
